# infra/terraform/files/rustc-perf-site.service
[Unit]
Description=rustc-perf site container
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
EnvironmentFile=/etc/rustc-perf-site.host.env
# The restore in ExecStartPre downloads and extracts the full backup archive
# on a fresh host, which takes well over systemd's default 90s start budget.
TimeoutStartSec=300
ExecStartPre=/usr/local/bin/rustc-perf-restore-if-empty
ExecStart=/usr/local/bin/rustc-perf-run-container
ExecStop=/usr/bin/docker stop rustc-perf-site
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
