# infra/terraform/files/rustc-perf-site.service
[Unit]
Description=rustc-perf site container
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# The restore in ExecStartPre downloads and extracts the full backup archive
# on a fresh host — minutes, and it grows with the database. A finite timeout
# would eventually kill a half-finished restore and retry it in a loop; the
# aws/tar pipeline inside the script already fails on its own on network
# errors, so this cannot hang silently.
TimeoutStartSec=infinity
ExecStartPre=/usr/local/bin/rustc-perf-restore-if-empty
ExecStart=/usr/local/bin/rustc-perf-run-container
ExecStop=/usr/bin/docker stop rustc-perf-site
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
